# services/orchestrator/systemd_template/ghostpaw-orchestrator.service
[Unit]
Description=Ghostpaw Orchestrator Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/melynxis/solace/services/orchestrator
ExecStart=/home/melynxis/solace/services/orchestrator/dev.sh
Restart=always
User=melynxis

[Install]
WantedBy=multi-user.target