# services/orchestrator/systemd_template/ghostpaw-mcp.service
[Unit]
Description=Ghostpaw MCP Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/melynxis/solace/services/orchestrator
ExecStart=/home/melynxis/solace/services/orchestrator/dev-mcp.sh
Restart=always
User=melynxis

[Install]
WantedBy=multi-user.target